# rental-car-checker-backend.service
[Unit]
Description=Rental Car Checker Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rental-car-checker
Environment="PATH=/home/ubuntu/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=production"
ExecStart=/home/ubuntu/.bun/bin/bun run backend/server.ts
Restart=always
RestartSec=10
StandardOutput=append:/var/log/rental-car-checker-backend.log
StandardError=append:/var/log/rental-car-checker-backend.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
